# etc/systemd/system/dragon-eye.service
[Unit]
Description=Dragon Eye
After=network.target nvargus-daemon.service

[Service]
ExecStart=/usr/local/bin/dragon-eye
Type=oneshot

[Install]
WantedBy=multi-user.target

